# flyerboard/deploy/flyerboard.service
[Unit]
Description=FlyerBoard Backend Service
After=network.target

[Service]
Type=simple
User=flyerboard
Group=flyerboard
WorkingDirectory=/opt/flyerboard

# Environment variables — copy env.example to /etc/flyerboard/env and fill in values
EnvironmentFile=/etc/flyerboard/env

ExecStart=/usr/bin/java -jar /opt/flyerboard/flyerboard-back-end.jar

Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/flyerboard

StandardOutput=journal
StandardError=journal
SyslogIdentifier=flyerboard

[Install]
WantedBy=multi-user.target
